# deploy/systemd/ktbot-discord-sender.service
[Unit]
Description=KTBot Discord digest sender
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/ktbot
Environment=PYTHONPATH=/home/ubuntu/ktbot/src
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/ktbot/.venv/bin/python -m discord_bot.send_digest
